# conf/yingsan.conf
[program:web20]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/web20
autostart=true
autorestart=true
